# tree/etc/services.d/metrics-server/run
#!/bin/bash
#

mkdir -p /var/www
touch /var/www/index.html

exec darkhttpd /var/www --port 9630 --no-listing --default-mimetype text/plain 

