# tree/etc/services.d/dhcp-relay/run
#!/bin/sh -e
#
config="/etc/dhcp/dhcrelay.conf"
source $config

exec dhcrelay -d $DHCP_DOWNSTREAM_INTERFACES $DHCP_UPSTREAM_INTERFACES $DHCP_UPSTREAM_ADDRESSES