# tree/etc/services.d/dmz-network/finish
#!/bin/sh
#
. /run/network/dmz-config

ip ad flush dev "$INTERFACE"
ip li set dev "$INTERFACE" state down
