# tree/etc/services.d/dmz-network/run
#!/bin/sh
. /run/network/dmz-config

ip li set dev "$INTERFACE" address "$HWADDR"

exec dhcpcd -B -C resolv.conf -h $DHCP_HOSTNAME -I $HWADDR $INTERFACE

