# tree/etc/services.d/confd/run
#!/bin/bash
#
/usr/local/bin/confd -watch
