# tree/etc/services.d/bird/run
#!/bin/bash
#

exec /usr/bin/bird -c /etc/bird.conf -f
